# playbooks/roles/django/templates/gunicorn.socket.j2
[Unit]
Description={{ project_settings.project_name }}-gunicorn socket

[Socket]
ListenStream=/run/{{ project_settings.project_name }}-gunicorn.sock

[Install]
WantedBy=sockets.target